# root/etc/crontabs/abc
* * * * *   /usr/bin/php /app/monica/artisan schedule:run
